# services/measure-dh.service.j2
[Unit]
Description=Measure Temperature and Humidity

[Service]
ExecStart={{ scripts.measure_dh.bin }}
WorkingDirectory={{ scripts.measure_dh.dir }}
Type=oneshot
