# services/piui.service.j2
[Unit]
Description=PIUI
Wants=nginx.service

[Service]
ExecStart={{ scripts.piui.dir }}/{{ scripts.piui.bin }}
WorkingDirectory={{ scripts.piui.dir }}
User=pi